# docker/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:saltmaster]
command=/usr/bin/salt-master -l info

[program:saltapi]
command=/usr/bin/salt-api

[program:approve-keys]
command = /usr/local/bin/approve-keys.sh
